# systemd/.config/systemd/user/ienabler.service
[Unit]
Description=Run ienabler


[Service]
ExecStart=%h/bin/pingtest
